# misc/supervisord.conf
[supervisord]
nodaemon=true

[program:httpd]
command=/usr/local/bin/http-server /app/public/ -p 80 -d false

[program:socket-server]
command=/usr/local/bin/npm start
directory=/app/
